# roles/docker.coreos.cleanup/files/docker-cleanup-dind.timer
[Unit]
Description=Run docker-cleanup-dind.service every hour

[Timer]
OnCalendar=hourly
