# roles/docker.coreos.cleanup/files/docker-cleanup.timer
[Unit]
Description=Run docker-cleanup.service every hour

[Timer]
OnCalendar=hourly
