# roles/docker.coreos.cleanup/files/docker-cleanup.service
[Unit]
Description=Runs /opt/bin/docker-cleanup.sh

[Service]
Type=oneshot
ExecStart=/opt/bin/docker-cleanup.sh
